# taleshed.service.example
[Unit]
Description=TaleShed MCP Server (HTTP)
After=network.target

[Service]
Type=simple
User=jpmoo
Group=jpmoo
WorkingDirectory=/home/jpmoo/taleshed
EnvironmentFile=-/home/jpmoo/taleshed/.env
ExecStart=/usr/bin/node dist/http.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
